# bitnet-mqtt-device.service
[Unit]
Description=BitNet MQTT Device Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/bitnet_runner2
Environment=PATH=/home/pi/bitnet_runner2/.venv/bin
ExecStart=/home/pi/bitnet_runner2/.venv/bin/python bitnet_mqtt_device.py --config config.json service
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
